# agent/deploy/com.bifrost.peer-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.bifrost.peer-watchdog</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-c</string>
    <string>source /Users/vision/bifrost-agent/config/env.sh &amp;&amp; source /Users/vision/bifrost-agent/config/.env 2>/dev/null; exec /Users/vision/bifrost-agent/peer_watchdog.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/vision/bifrost-agent</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/vision</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>StartInterval</key>
  <integer>60</integer>

  <key>StandardOutPath</key>
  <string>/Users/vision/bifrost-agent/logs/peer-watchdog-launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/vision/bifrost-agent/logs/peer-watchdog-launchd.err</string>
</dict>
</plist>
