# agent/schedules/com.bifrost.remediation-runner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.bifrost.remediation-runner</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>cd /Users/vision/Desktop/stocks/bifrost-platform &amp;&amp; python3 scripts/run_remediation_agent.py --no-watch --install</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/vision/Desktop/stocks/bifrost-platform/agent/reports/runner-launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/vision/Desktop/stocks/bifrost-platform/agent/reports/runner-launchd.err</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
