# agent/git-bridge/com.bifrost.git-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.bifrost.git-bridge</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/vision-mac-trader/Desktop/stocks/bifrost-platform/agent/git-bridge/run.sh</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/vision-mac-trader</string>
    <key>GIT_BRIDGE_PORT</key>
    <string>8785</string>
    <key>GIT_WORKSPACE_ROOT</key>
    <string>/Users/vision-mac-trader/Desktop/stocks</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/vision-mac-trader/Desktop/stocks/bifrost-platform/agent/git-bridge/logs/git-bridge.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/vision-mac-trader/Desktop/stocks/bifrost-platform/agent/git-bridge/logs/git-bridge.err</string>

  <key>ThrottleInterval</key>
  <integer>10</integer>
</dict>
</plist>
